# molecule/default/config/testservice/root_test.timer
# test.timer

[Unit]
Description=Root Test Timer

[Timer]
OnCalendar=*:*:0
Unit=root_test.service

[Install]
WantedBy=timers.target
